# cronjobs/send_logs.cron
# this should be added to the cron for user "uglymug"
# you'll need to change the path to the script
# I'd recommend setting the MAILTO variable for the
# user, if you don't already

# MIN HOUR DAY MONTH DAYOFWEEK   COMMAND
# at 3:10am every day - 10 minutes should be long enough for kill_scat.sh
# to work it's mojo
10 3 * * * /home/chisel/Development/SourceForge/UglyLogs/cronjobs/send_logs.sh
